# packages/utils/crond/config-files/etc/cron.weekly/vacuumsqlite
#!/bin/sh

sqlite3 /etc/params/wymedia/.wyplay_db.-1.db 'VACUUM;'
