# packages/utils/crond/config-files/etc/cron.daily/ntp
#!/bin/sh

NTPLOG=/wymedia/usr/var/log/ntpdate.log

if [ -f /etc/ntp.conf ]; then
        NTPSERVER=`grep "server" /etc/ntp.conf |head -n1 |cut -d " " -f2`
else
        NTPSERVER="0.pool.ntp.org"
fi

ntpdate $NTPSERVER >>$NTPLOG
